# license-server/deploy/ocorrencias-worker.service
[Unit]
Description=Processador de Ocorrencias — worker RQ
After=network.target redis-server.service

[Service]
User=ubuntu
WorkingDirectory=/home/ubuntu/license-server
EnvironmentFile=/home/ubuntu/license-server/.env
ExecStart=/home/ubuntu/license-server/.venv/bin/rq worker default
Restart=on-failure

[Install]
WantedBy=multi-user.target
